# container-addons/debian-9-lxde/home/user/.config/autostart/sd-user-workaround-xauthority.desktop
[Desktop Entry]
Name=Workaround for systemd dbus XAUTHORITY
Exec=dbus-update-activation-environment --systemd XAUTHORITY
Terminal=false
Type=Application
