# container-addons/debian-9-lxde/home/user/.config/autostart/sd-user-workaround-display.desktop
[Desktop Entry]
Name=Workaround for systemd dbus DISPLAY
Exec=dbus-update-activation-environment --systemd DISPLAY
Terminal=false
Type=Application
